# deploy/prediction-pulse-prune.timer
[Unit]
Description=Daily prediction-pulse snapshot prune

[Timer]
# Once a day, early-morning local time. Persistent catches up a missed run after downtime.
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
